# overlay/etc/services.d/reload-listener/run
#!/usr/bin/with-contenv sh

set -euo pipefail
export BF_E="${PWD##*/}/$(basename ${0})"


#======================================================================================================================
# Listen for .
#======================================================================================================================

bf-echo "Starting TCP Listener."
s6-tcpserver4 -v 0.0.0.0 5000 postfix-reload-listener
